# website/control/brp.baseline.test.cron
0 1 * * * /opt/mrs-biomarkers/pipeline/brp.baseline.test
